# BSD-exporter/prometheus_exporter_start
#!/bin/sh
. /etc/rc.subr

SCRIPT_PATH="/usr/local/bin/prometheus_exporter"
bash "$SCRIPT_PATH"
